# root/etc/services.d/ssmgr-webui/run
#!/usr/bin/with-contenv sh

if [ -s /root/.ssmgr/webgui.yml ]; then
    cd /ssmgr
    exec node server.js -c /root/.ssmgr/webgui.yml
fi