# root/etc/services.d/ssmgr/run
#!/usr/bin/with-contenv sh
cd /ssmgr
exec node server.js -c /root/ss.yml